# salt/circus/circus.conf
# circus process manager

start on runlevel 2
start on runlevel 3
start on runlevel 4
start on runlevel 5

stop on shutdown

respawn
exec /usr/local/bin/circusd --log-output /var/log/circus.log \
                      --pidfile /var/run/circusd.pid \
                      /etc/circus.ini